# Week10-SmallAddressMaster/AMaven.service
[Unit]
Description=Address Maven

[Service]
ExecStart=/home/ubuntu/npm/bin/npm start
<<<<<<< HEAD
WorkingDirectory=/home/ubuntu/bin/address-maven
=======
WorkingDirectory=/home/ubuntu/bin/AMaven
>>>>>>> 8fa2ce617497c7d06ffcefd3bb4b876c8e98d738
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=AMaven
User=ubuntu
Group=ubuntu
Environment=NODE_ENV=production
Environment=PORT=30029

[Install]
WantedBy=multi-user.target
